# sysmotd/sysmotd.timer
[Unit]
Description=Timer for sysmotd service

[Timer]
Unit=sysmotd.service
OnActiveSec=15s
OnUnitActiveSec=1m  

[Install]
WantedBy=timers.target