# rpmbuild/BUILD/sysmotd/sysmotd.timer
[Unit]
Description=Timer for sysmotd service
Requires=sysmotd.service

[Timer]
Unit=sysmotd.service
OnBootSec=15s
OnUnitActiveSec=1m  

[Install]
WantedBy=timers.target